# stage3/02-install-eric/files/eric.service
[Unit]
Description=Run Eric
[Service]
ExecStart=/bin/sh -c "sudo -H -u pi /home/pi/eric/main"
[Install]
WantedBy=multi-user.target
