# stage3/03-install-blas/files/blas-supervisord.conf
[program:blas]
command = gunicorn -w 9 -b 0.0.0.0:5000 main:app
directory = /home/pi/blas
autostart = true
autorestart = true
